# crontab
* * * * * root /usr/local/bin/python3 /app/parse.py >> /var/log/cron.log 2>&1
